# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:vnc]
command=/bin/bash -c "su - desktop -c 'vncserver :1 -geometry ${VNC_RESOLUTION:-1366x768} -depth 24 -localhost no'"
priority=10
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --listen %(ENV_PORT)s --vnc localhost:5901
priority=20
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:pulseaudio]
command=/bin/bash -c "su - desktop -c 'pulseaudio --start --exit-idle-time=-1'"
priority=5
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
